# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:haumea]
directory=/data/www/haumea
command=/usr/local/bin/gunicorn -c /etc/default/gunicorn haumea.wsgi